# vps/setup/opencode.service
[Unit]
Description=OpenCode headless server
After=network.target

[Service]
Type=simple
User=opencode
Group=opencode
# Without this, systemd's ProtectHome=read-only blocked writes under ~/ and OpenCode exited immediately.
Environment=HOME=/opt/leadagent
WorkingDirectory=/opt/leadagent
EnvironmentFile=/etc/opencode.env
Environment=XDG_DATA_HOME=/opt/leadagent/data
Environment=XDG_CONFIG_HOME=/opt/leadagent/config
Environment=XDG_CACHE_HOME=/opt/leadagent/cache
Environment=OPENCODE_DISABLE_AUTOUPDATE=true
Environment=OPENCODE_DISABLE_LSP_DOWNLOAD=true
Environment=OPENCODE_DISABLE_CLAUDE_CODE=true
Environment=OPENCODE_DISABLE_TERMINAL_TITLE=true
Environment=OPENCODE_AUTO_SHARE=false
ExecStart=/usr/local/bin/opencode serve --hostname 127.0.0.1 --port 4096
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# HOME must be under ReadWritePaths: strict ProtectHome prevented runtime cache writes.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=yes
ProtectHome=no
ReadWritePaths=/opt/leadagent

[Install]
WantedBy=multi-user.target
